# docker/php-fpm/crontab
# Github subscriptions
0 * * * * php /app/src/index.php update-github-subscriptions

